# supervisor/conf.d/coffee.conf
[program:coffee]
command=/home/benjamin/coffee/venv/bin/gunicorn -c gunicorn_config.py wsgi:application
directory=/home/benjamin/coffee/coffee
user=benjamin