# supervisor/conf.d/hello.conf
[program:hello]
command=/home/benjamin/env/gunicorn/bin/gunicorn --bind 0.0.0.0:9081 -k gevent wsgi:app
directory=/home/benjamin/env
#user=benjamin
